# deploy/systemd/reputation-oracle@.service
[Unit]
Description=Kleros Reputation Oracle (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=no
User=oracle
Group=oracle
WorkingDirectory=/opt/reputation-oracle/bot
EnvironmentFile=/etc/reputation-oracle/%i.env
ExecStart=/usr/bin/node --import tsx /opt/reputation-oracle/bot/src/index.ts
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reputation-oracle-%i
TimeoutStartSec=300
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
